# upstart.conf
# ghost

# description "An Upstart task to make sure that my Ghost server is always running"
# author "Leonardo Iacovini"

start on startup

script
    cd ~/blog
    npm start --production
end script